# roles/swhkd/files/.config/systemd/user/swhkd.service
[Unit]
Description=swhkd hotkey daemon
BindsTo=default.target

[Service]
Type=simple
ExecStart=%h/sh/hotkeys.sh

[Install]
WantedBy=default.target
