# roles/sxhkd/files/.config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Comment=Simple X hotkey daemon
Exec=/usr/bin/sxhkd 2>&1 | tee -a /dev/shm/sxhkd.log
Terminal=false
Type=Application
